# deploy/linux/nixis-daemon.service
[Unit]
Description=Nixis Governance Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/nixis-daemon -policy-dir /etc/nixis/policies -audit-db /var/lib/nixis/audit.db
Restart=on-failure
RestartSec=5
User=nixis
Group=nixis
Environment=NIXIS_SOCKET_PATH=/tmp/nixis.sock
StandardOutput=journal
StandardError=journal
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
