# deploy/macos/com.aegis.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aegis.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/aegis-daemon</string>
        <string>-policy-dir</string>
        <string>/etc/aegis/policies</string>
        <string>-audit-db</string>
        <string>/var/lib/aegis/audit.db</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>StandardOutPath</key>
    <string>/var/log/aegis/daemon.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/aegis/daemon.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>AEGIS_SOCKET_PATH</key>
        <string>/tmp/aegis.sock</string>
    </dict>
</dict>
</plist>
