# policies/imported/falco/falco-schedule-cron-jobs.yaml
# IMPORT_REVIEW: translated by Claude — verify semantics
# imported from: rules via aegis policy import
apiVersion: aegis.io/v1
kind: PolicyTemplate
metadata:
    name: falco-schedule-cron-jobs
    annotations:
        aegis.io/falco-tags: maturity_incubating,host,container,filesystem,mitre_execution,T1053.003
        aegis.io/imported-from: aegis-import-2621062960.yaml
        aegis.io/severity: low
        aegis.io/llm-confidence: medium
        aegis.io/source-rule: Schedule Cron Jobs
spec:
    description: |
        Detect scheduled cron jobs; this is a highly generic detection and certainly needs adjustments and profiling in your environment before operationalization. Simultaneously, exploiting the functionality of cron jobs is among one of the oldest TTPs used by adversaries.
    matchConstraints:
        tools:
            - Bash
            - Write
            - Edit
    validations:
        - expression: 'tool == "Bash" && (request.args.command.contains("crontab") || request.args.command.matches(".*>\\s*/etc/cron\\.(d|daily|hourly|weekly|monthly)/.*")) || (tool == "Write" || tool == "Edit") && request.args.path.matches(".*/etc/cron.*")'
          message: 'Falco rule violated: Schedule Cron Jobs'
          action: REQUIRE_APPROVAL
    defaultAction: ALLOW
