# policies/imported/falco/falco-schedule-cron-jobs.yaml
# IMPORT_TODO: Falco condition uses kernel or macro fields that cannot be automatically translated — manual review required. Original condition: ((open_write and fd.name startswith /etc/cron) or
#  (spawned_process and proc.name = "crontab"))
# and not user_known_cron_...
# imported from: main via aegis policy import
apiVersion: aegis.io/v1
kind: PolicyTemplate
metadata:
    name: falco-schedule-cron-jobs
    annotations:
        aegis.io/falco-tags: maturity_incubating,host,container,filesystem,mitre_execution,T1053.003
        aegis.io/imported-from: aegis-import-487448035.yaml
        aegis.io/severity: low
        aegis.io/source-rule: Schedule Cron Jobs
spec:
    description: |
        Detect scheduled cron jobs; this is a highly generic detection and certainly needs adjustments and profiling in your environment before operationalization. Simultaneously, exploiting the functionality of cron jobs is among one of the oldest TTPs used by adversaries.
    matchConstraints:
        tools:
            - Read
            - Write
            - Edit
    validations:
        - expression: "false"
          message: 'Falco rule violated: Schedule Cron Jobs'
          action: REQUIRE_APPROVAL
    defaultAction: ALLOW
